# features/epel/epel-sync.service
[Unit]
Description=Synchronize and index mail with epel

[Service]
Type=oneshot
ExecStart=%h/.local/bin/epel sync
